# run/fangli.service
[Service]
ExecStart=/opt/tiger/fangli/run/real_run
Restart=always
[Install]
WantedBy=default.target